# systemd/maintain-lanhosts-route.target
[Unit]
Description=Chain of the maintainer for lanhosts route table
RefuseManualStart=no

